# roles/docker/files/etc/profile.d/compose.sh
#!/usr/bin/env bash

export PATH="/usr/libexec/docker/cli-plugins:$PATH"
