# roles/login/templates/etc/systemd/system/getty@tty1.service.d/autologin.conf.j2
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin {{ user.name }} %I $TERM
Type=simple
Environment=XDG_SESSION_TYPE=x11